# dev/supervisord-sidekiq.conf
[supervisord]
nodaemon=true
user=root

[program:sidekiq]
command=bundle exec sidekiq -e production
directory=/app
autostart=true
autorestart=false
stopasgroup=true
stopsignal=TERM