# dev/supervisord-puma.conf
[supervisord]
nodaemon=true
user=root

[program:puma]
command=bundle exec puma -b tcp://0.0.0.0 -p 3000
directory=/app
autostart=true
autorestart=false
stopasgroup=true
stopsignal=TERM