# dev/supervisord-sneakers.conf
[supervisord]
nodaemon=true
user=root

[program:sneakers]
command=bundle exec rake sneakers:run
directory=/app
autostart=true
autorestart=false
stopasgroup=true
stopsignal=TERM