# dev/supervisord-whenever.conf
[supervisord]
nodaemon=true
user=root

[program:whenever]
command=/bin/sh ./dev/whenever.sh
directory=/app
autostart=true
autorestart=false
stopasgroup=true
stopsignal=TERM