# dev/supervisord-clockwork.conf
[supervisord]
nodaemon=true
user=root

[program:clockwork]
command=bundle exec clockwork clock.rb
directory=/app
autostart=true
autorestart=false
stopasgroup=true
stopsignal=TERM